# scripts/podman/remoteclaw.container.in
# RemoteClaw gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into remoteclaw's ~/.config/containers/systemd/
# {{REMOTECLAW_HOME}} is replaced at install time.

[Unit]
Description=RemoteClaw gateway (rootless Podman)

[Container]
Image=remoteclaw:local
ContainerName=remoteclaw
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{REMOTECLAW_HOME}}/.remoteclaw:/home/node/.remoteclaw
EnvironmentFile={{REMOTECLAW_HOME}}/.remoteclaw/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
